# supervisord.conf
[supervisord]
nodaemon=true

[program:celery-worker]
command=celery -A app.celery_config:celery_app worker --loglevel=info -E --concurrency=2
autostart=true
autorestart=true
# stderr_logfile=/var/log/celery-worker.err.log
# stdout_logfile=/var/log/celery-worker.out.log
stderr_logfile=/dev/stdout
stdout_logfile=/dev/stdout

[program:celery-beat]
command=celery -A app.celery_config:celery_app beat --loglevel=info
autostart=true
autorestart=true
# stderr_logfile=/var/log/celery-beat.err.log
# stdout_logfile=/var/log/celery-beat.out.log
stderr_logfile=/dev/stdout
stdout_logfile=/dev/stdout